# scripts/supervisord-craftbukkit.conf
[program:craftbukkit]
command=/opt/minecraft/run.sh
numprocs=1
autostart=true
autorestart=true
